# imooc/conf/supervisor_celery_worker.ini
[program:imooc-celery-worker]
command=python manage.py celery worker -l INFO
directory=/Users/ych/devs/workspace/yidianzhishi/course/celery-learning/imooc
environment=PATH="/Users/ych/.pyenv/versions/env4mooc/bin"
stdout_logfile=/Users/ych/devs/workspace/yidianzhishi/course/celery-learning/imooc/logs/celery.worker.log
stderr_logfile=/Users/ych/devs/workspace/yidianzhishi/course/celery-learning/imooc/logs/celery.worker.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=997
