# imooc/conf/supervisor_celery_flower.ini
[program:imooc-celery-flower]
command=python manage.py celery flower
directory=/Users/ych/devs/workspace/yidianzhishi/course/celery-learning/imooc
environment=PATH="/Users/ych/.pyenv/versions/env4mooc/bin"
stdout_logfile=/Users/ych/devs/workspace/yidianzhishi/course/celery-learning/imooc/logs/celery.flower.log
stderr_logfile=/Users/ych/devs/workspace/yidianzhishi/course/celery-learning/imooc/logs/celery.flower.log
autostart=true
autorestart=true
startsecs=10
stopwatisecs=60
priority=995
